# root/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

echo "Starting cron"

exec cron -f
